# internal/assets/docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Xvnc: virtual display + VNC server ────────────────────────────────────────
[program:xvnc]
command=Xvnc :1 -geometry 1280x800 -depth 24 -SecurityTypes None -localhost no
priority=10
autostart=true
autorestart=true
user=node
environment=HOME="/home/node",USER="node"
stdout_logfile=/tmp/xvnc.log
stderr_logfile=/tmp/xvnc.err.log

; ── XFCE4 desktop session ──────────────────────────────────────────────────────
[program:xfce4]
command=bash -c "sleep 3 && startxfce4"
priority=20
autostart=true
autorestart=unexpected
user=node
environment=DISPLAY=":1",HOME="/home/node",USER="node",DBUS_SESSION_BUS_ADDRESS="autolaunch:"
stdout_logfile=/tmp/xfce4.log
stderr_logfile=/tmp/xfce4.err.log

; ── noVNC: browser-accessible VNC ─────────────────────────────────────────────
[program:novnc]
command=websockify --web /usr/share/novnc 0.0.0.0:6901 localhost:5901
priority=30
autostart=true
autorestart=true
user=node
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.err.log

; ── OpenClaw Gateway ──────────────────────────────────────────────────────────
; autostart=false: the gateway cannot start until 'openclaw onboard' has been run.
; After onboarding, ClawSandbox starts it via: supervisorctl start openclaw
; autorestart=true: if the gateway crashes after successful startup, supervisord
; will restart it automatically. This does NOT conflict with supervisorctl stop —
; manually stopped processes stay stopped regardless of autorestart setting.
[program:openclaw]
command=openclaw gateway --port 18789
priority=40
autostart=false
autorestart=true
user=node
directory=/home/node
environment=HOME="/home/node",USER="node",DISPLAY=":1",PATH="/usr/local/bin:/usr/bin:/bin",PLAYWRIGHT_BROWSERS_PATH="/ms-playwright"
stdout_logfile=/tmp/openclaw.log
stderr_logfile=/tmp/openclaw.err.log
